# rootfs/etc/services.d/01-confd/run
#!/usr/bin/with-contenv sh

exec /opt/anaxexp/tools/bin/confd -watch
